# setup-env/citytouch.conf
[program:citytouch]
user=root
command=/data/scripts/start-citytouch.sh
autorestart=true
stopsignal=INT